# com.deepthinktrader.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.deepthinktrader.bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/rebelhawk/Documents/Claude/StockTrader/main.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/rebelhawk/Documents/Claude/StockTrader</string>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/rebelhawk/Documents/Claude/StockTrader/deepthinktrader.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/rebelhawk/Documents/Claude/StockTrader/deepthinktrader_error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/rebelhawk/Library/Python/3.9/bin</string>
    </dict>
</dict>
</plist>
